# constructr-coordination-demo/docker/supervisord.conf
[supervisord]
http_port=/opt/supervisor/supervisor.sock ; (default is to run a UNIX domain socket server)
logfile=/opt/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=/opt/supervisor/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[supervisorctl]

[group:app]
programs=consul,constructr-consul-demo

[program:consul]
directory=/opt/consul
command=/bin/consul agent -client=0.0.0.0 -join=%(ENV_CONSUL_SERVER_ADDRESS)s -data-dir=/opt/consul
stdout_logfile=/opt/consul/logs/stdout.log
stderr_logfile=/opt/consul/logs/error.log
numprocs=1
autostart=true
autorestart=true
stopsignal=INT
priority=900

[program:constructr-consul-demo]
directory=/mnt
command=/mnt/bin/constructr-coordination-demo
stdout_logfile=/opt/constructr-coordination-demo/logs/stdout.log
stderr_logfile=/opt/constructr-coordination-demo/logs/error.log
numprocs=1
autostart=true
autorestart=true
stopsignal=INT
priority=901
